# bridge-box.service
[Unit]
Description=Bridge Scoring Box Startup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/bridgebox/bridge-box/startup.sh
RemainAfterExit=yes
User=bridgebox
WorkingDirectory=/home/bridgebox

[Install]
WantedBy=multi-user.target
